# deploy/roles/gunicorn/templates/systemd.socket.j2
[Unit]
Description={{ gunicorn_app_name }} socket

[Socket]
ListenStream=/run/{{ gunicorn_app_name }}.sock
SocketUser={{ nginx_user }}
SocketMode=600

[Install]
WantedBy=sockets.target
