# deploy/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock
SocketUser=nginx
SocketMode=600

[Install]
WantedBy=sockets.target
